# src/systemd/cockpit-wsinstance-https@.service.in
[Unit]
Description=Cockpit Web Service https instance %I
BindsTo=cockpit.service
Documentation=man:cockpit-ws(8)
Requires=cockpit-session.socket cockpit-session-socket-owner.service
After=cockpit-session.socket cockpit-session-socket-owner.service

[Service]
Slice=system-cockpithttps.slice
ExecStart=@libexecdir@/cockpit-ws --for-tls-proxy --port=0
DynamicUser=yes
SupplementaryGroups=cockpit-session-socket
