# src/systemd/cockpit-wsinstance-http.socket.in
[Unit]
Description=Socket for Cockpit Web Service http instance
BindsTo=cockpit.service
Documentation=man:cockpit-ws(8)

[Socket]
ListenStream=/run/cockpit/wsinstance/http.sock
SocketUser=cockpit-ws
SocketMode=0600
RemoveOnStop=yes
