# src/systemd/cockpit-session@.service.in
[Unit]
Description=Cockpit session %I

[Service]
ExecStart=@libexecdir@/cockpit-session
StandardInput=socket
StandardOutput=inherit
StandardError=journal
User=root
