# src/systemd/cockpit-wsinstance-https-factory@.service.in
[Unit]
Description=Cockpit Web Service https instance factory
Documentation=man:cockpit-ws(8)

[Service]
ExecStart=@libexecdir@/cockpit-wsinstance-factory
User=root
