# src/systemd/cockpit-wsinstance-https-factory.socket.in
[Unit]
Description=Socket for Cockpit Web Service https instance factory
BindsTo=cockpit.service
Documentation=man:cockpit-ws(8)

[Socket]
ListenStream=/run/cockpit/wsinstance/https-factory.sock
Accept=yes
SocketUser=cockpit-ws
SocketMode=0600
RemoveOnStop=yes
