# src/systemd/cockpit-motd.service.in
[Unit]
Description=Cockpit motd updater service
Documentation=man:cockpit-ws(8)
Wants=network-online.target
After=network-online.target cockpit.socket

[Service]
Type=oneshot
ExecStart=-@datadir@/@PACKAGE@/motd/update-motd
