# src/systemd/cockpit-wsinstance-http.service.in
[Unit]
Description=Cockpit Web Service http instance
Documentation=man:cockpit-ws(8)
BindsTo=cockpit.service
Requires=cockpit-session.socket cockpit-session-socket-user.service
After=cockpit-session.socket cockpit-session-socket-user.service

[Service]
ExecStart=@libexecdir@/cockpit-ws --no-tls --port=0
User=cockpit-wsinstance
Group=cockpit-wsinstance
SupplementaryGroups=cockpit-session-socket
