# data/cockpit-ws-notls.service.in
[Unit]
Description=Cockpit Web Server (No TLS)
Documentation=man:cockpit-ws(8)
Requires=cockpit-ws-notls.socket
Conflicts=cockpit-ws.service

[Service]
ExecStart=@libexecdir@/cockpit-ws --no-tls
