# data/test-server.service.in
[Unit]
Description=Test Server

[Service]
WorkingDirectory=@datadir@/cockpit-test-assets
ExecStart=@datadir@/cockpit-test-assets/test-server
