# data/test-server.socket
[Unit]
Description=Test Server Socket

[Socket]
ListenStream=8765
