# src/ws/cockpit-wsinstance-https.socket.in
[Unit]
Description=Socket for Cockpit Web Service https instance
BindsTo=cockpit.service
Documentation=man:cockpit-ws(8)

[Socket]
ListenStream=/run/cockpit/wsinstance/https.sock
SocketUser=@user@
SocketMode=0600
