# src/ws/cockpit.service.in
[Unit]
Description=Cockpit Web Service
Documentation=man:cockpit-ws(8)
Requires=cockpit.socket

[Service]
RuntimeDirectory=cockpit/tls
RuntimeDirectoryMode=0700
# systemd ≥ 241 sets this automatically
Environment=RUNTIME_DIRECTORY=/run/cockpit/tls
ExecStartPre=@sbindir@/remotectl certificate --ensure --user=root --group=@group@ --selinux-type=@selinux_config_type@
ExecStart=@libexecdir@/cockpit-tls
PermissionsStartOnly=true
User=@user@
Group=@group@
