# src/ws/cockpit-wsinstance-http.service.in
[Unit]
Description=Cockpit Web Service http instance
BindsTo=cockpit.service
Documentation=man:cockpit-ws(8)

[Service]
ExecStart=@abs_libexecdir@/cockpit-ws --no-tls --port=0
User=@cockpit_ws_instance_user@
Group=@cockpit_ws_instance_group@
