# src/ws/cockpit-testing.service.in
[Unit]
Description=Cockpit Web Server (Testing)
Documentation=man:cockpit-ws(8)
Requires=cockpit-testing.socket
Conflicts=cockpit.service

[Service]
Environment=G_MESSAGES_DEBUG=cockpit-ws G_DEBUG=fatal-criticals G_SLICE=always-malloc,debug-blocks
ExecStart=@libexecdir@/cockpit-ws --no-tls
User=@user@
Group=@group@
