# src/ws/cockpit-testing.service.in
[Unit]
Description=Cockpit Web Server (Testing)
Documentation=man:cockpit-ws(8)
Requires=cockpit-testing.socket
Conflicts=cockpit.service

[Service]
Environment=G_MESSAGES_DEBUG=cockpit-ws,cockpit-bridge G_DEBUG=fatal-criticals G_SLICE=always-malloc,debug-blocks
ExecStartPre=-/bin/sh -c 'echo 0 > /proc/sys/kernel/yama/ptrace_scope'
ExecStart=@libexecdir@/cockpit-ws --no-tls
PermissionsStartOnly=true
User=@user@
Group=@group@
