# src/ws/cockpit-wsinstance-https@.service.in
[Unit]
Description=Cockpit Web Service https instance %I
BindsTo=cockpit.service
Documentation=man:cockpit-ws(8)

[Service]
ExecStart=@libexecdir@/cockpit-ws --for-tls-proxy --port=0
User=@wsinstanceuser@
Group=@wsinstancegroup@
