# src/ws/cockpit.service.in
[Unit]
Description=Cockpit Web Service
Documentation=man:cockpit-ws(8)
Requires=cockpit.socket

[Service]
ExecStartPre=@sbindir@/remotectl certificate --ensure --user=root --group=@group@ --selinux-type=@selinux_config_type@
ExecStart=@libexecdir@/cockpit-ws
PermissionsStartOnly=true
User=@user@
Group=@group@
