# src/ws/cockpit.service.in
[Unit]
Description=Cockpit Web Server
Documentation=man:cockpit-ws(8)
Requires=cockpit.socket

[Service]
ExecStart=@libexecdir@/cockpit-ws
