# src/ws/test-server.service.in
[Unit]
Description=Test Server

[Service]
WorkingDirectory=@datadir@/cockpit-test-assets
ExecStart=/usr/bin/runcon -t unconfined_service_t -- @datadir@/cockpit-test-assets/test-server
